# debian/intelmq-manager.prerm
#!/bin/sh
set -e

if [ "$1" = "remove" ] ; then
  a2dissite intelmq
  invoke-rc.d apache2 reload
fi

#DEBHELPER#
